# systemd/knot-resolver.service.in
[Unit]
Description=Knot Resolver Manager
Documentation=man:knot-resolver.systemd(7)
Wants=network-online.target
After=network-online.target
Before=nss-lookup.target
Wants=nss-lookup.target

[Service]
Type=notify
TimeoutStartSec=600s
ExecStart=@bin_dir@/knot-resolver --config=@etc_dir@/config.yaml
Environment=KRES_LOGGING_TARGET=syslog
ExecReload=@bin_dir@/kresctl --config=@etc_dir@/config.yaml reload
KillSignal=SIGINT
WorkingDirectory=@systemd_work_dir@
User=@user@
Group=@group@
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SETPCAP
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_SETPCAP
OOMPolicy=continue

[Install]
WantedBy=multi-user.target
