# systemd/kres-cache-gc.service.in
# SPDX-License-Identifier: CC0-1.0
[Unit]
Description=Knot Resolver Garbage Collector daemon
Documentation=man:kresd.systemd(7)
Documentation=man:kresd(8)

[Service]
Type=simple
ExecStart=@sbin_dir@/kres-cache-gc -c @systemd_cache_dir@ -d 1000
User=@user@
Group=@group@
Restart=on-failure
RestartSec=30
StartLimitInterval=400
StartLimitBurst=10
Slice=system-kresd.slice

[Install]
WantedBy=kresd.target
