# systemd/kres-watcher.service.in
[Unit]
Description=Knot Resolver Watcher daemon
Documentation=man:kresd.systemd(7)
Documentation=man:kresd(8)

[Service]
Type=simple
ExecStart=@sbin_dir@/kres-watcher
User=root
Group=root
Restart=on-failure
RestartSec=5
StartLimitInterval=400
StartLimitBurst=10
Slice=system-kresd.slice

[Install]
WantedBy=kresd.target
