# systemd/kresd-control@.socket.in
[Unit]
Description=Knot Resolver control socket
Documentation=man:kresd.systemd(7)
Documentation=man:kresd(8)
Before=sockets.target
PartOf=kresd@%i.service

[Socket]
Service=kresd@%i.service
ListenStream=@run_dir@/control@%i
FileDescriptorName=control
SocketMode=0660
Slice=system-kresd.slice

[Install]
WantedBy=sockets.target
