# systemd/nosocket/kres-cache-gc.service.in
[Unit]
Description=Knot Resolver Garbage Collector daemon
Documentation=man:kresd.systemd.nosocket(7)
Documentation=man:kresd(8)

[Service]
Type=simple
ExecStart=@sbin_dir@/kres-cache-gc -c @systemd_work_dir@ -d 1000
User=knot-resolver
Restart=on-abnormal
Slice=system-kresd.slice

[Install]
WantedBy=kresd.target
