# manager/knot_resolver_manager/kresd_controller/supervisord/supervisord.conf.j2
[supervisord]
pidfile = {{ config.pid_file }}
directory = /tmp
nodaemon = false
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
user=root

[unix_http_server]
file = {{ config.unix_http_server }}

[supervisorctl]
serverurl = unix://{{ config.unix_http_server }}

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface




{% for instance in config.instances %}

[program:{{ instance.id }}]

{%- if instance.type.name == "KRESD" %}

directory=/var/lib/knot-resolver
command=/usr/sbin/kresd -c /usr/lib/knot-resolver/distro-preconfig.lua -c /etc/knot-resolver/kresd.conf -n
environment=SYSTEMD_INSTANCE={{ instance.id }}
{%- elif instance.type.name == "GC" %}

directory=/var/lib/knot-resolver
command=/usr/sbin/kres-cache-gc -c /var/cache/knot-resolver -d 1000

{%- else %}

{# other subprocess types are not implemented, fail #}
{{ 0 / 0 }}

{% endif %}

{%- endfor -%}