# manager/knot_resolver_manager/kresd_controller/supervisord/supervisord.conf.j2
[supervisord]
pidfile = {{ config.pid_file }}
directory = {{ config.workdir }}
nodaemon = false
logfile = {{ config.logfile }}
logfile_maxbytes = 50MB
loglevel = trace
{# user=root #}

[unix_http_server]
file = {{ config.unix_http_server }}

[supervisorctl]
serverurl = unix://{{ config.unix_http_server }}

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:notify]
supervisor.rpcinterface_factory = knot_resolver_manager.kresd_controller.supervisord.plugin:make_rpcinterface


[program:kresd]
process_name=%(program_name)s%(process_num)d
numprocs={{ kresd.max_procs }}
redirect_stderr=false
stdout_logfile={{ kresd.logfile }}
stderr_logfile={{ kresd.logfile }}
directory={{ kresd.workdir }}
command={{ kresd.command }}
autostart=false
autorestart=true
stopsignal=TERM
killasgroup=true
startsecs=1
environment={{ kresd.environment }}

[program:gc]
redirect_stderr=false
stdout_logfile={{ gc.logfile }}
stderr_logfile={{ gc.logfile }}
directory={{ gc.workdir }}
command={{ gc.command }}
autostart=false
autorestart=true
stopsignal=TERM
killasgroup=true
startsecs=1
environment={{ gc.environment }}