# systemd/kresd@.service.d/module-http.conf
[Service]
Sockets=kresd-doh.socket
Sockets=kresd-webmgmt.socket
